# IIoTDB/db_iiot.service
[Unit]
Description=DBIIoT

[Service]
ExecStart=/opt/db_iiot/start.sh

[Install]
WantedBy=multi-user.target