# IIoTWeb/web_iiot.service
[Unit]
Description=WebService

[Service]
ExecStart=/opt/web_service/start.sh

[Install]
WantedBy=multi-user.target
